# conf.d/mongod.conf
[program:mongod]
command=mongod
autorestart=true
