# conf.d/rails.conf
[program:rails]
command=bundle exec rails server -p 80 -b 0.0.0.0
autorestart=true
