# conf.d/tornado.conf
[program:tornado]
command=python3 main.py
directory=/var/www
autorestart=true
